# deploy/systemd/station-converter-ja-ingest-n02.service
[Unit]
Description=station_converter_ja N02 ingest
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/station_converter_ja
User=station-converter-ja
Group=station-converter-ja
EnvironmentFile=/etc/station_converter_ja/station.env
RuntimeDirectory=station-converter-ja
ExecStart=/opt/station_converter_ja/target/release/station-ops job ingest-n02 $STATION_INGEST_ARGS
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
